# meta-nexell-distro/recipes-extended/nexell-nxupdate-service/files/nx_update/lib/systemd/system/nxupdate.service
[Unit]
Description=Nexell update
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/etc/update.sh

[Install]
WantedBy=multi-user.target
