# meta-nxp3220/recipes-bbappend/android-tools/files/android-tools-adbd.service
[Unit]
Description=Android Debug Bridge daemon

[Service]
Type=forking
ExecStart=/usr/bin/start_adbd.sh
PIDFile=/run/adbd.pid

[Install]
WantedBy=multi-user.target
