# recipes-multimedia/amixer-conf/files/amixer-conf.service
[Unit]
Description=amixer conf

[Service]
Type=oneshot
ExecStart=/usr/bin/amixer-conf.sh

[Install]
WantedBy=multi-user.target