# meta-nxp3220/recipes-bbappend/swupdate/files/swupdate-progress.service
[Unit]
Description=swupdate progress service

[Service]
ExecStart=/usr/bin/swupdate-progress -r -w

[Install]
WantedBy=multi-user.target
