# meta-nexell-distro/recipes-qt/qt5_5.8.x/qtcinema-launcher/qt_cinema.service
[Unit]
Description=QT5 cinema demo

[Service]
Type=oneshot
ExecStart=/usr/bin/qt_cinema_run.sh

[Install]
WantedBy=multi-user.target
