# meta-nexell-distro/recipes-extended/nexell-nxupdate/files/32/lib/systemd/system/nxupdate.service
[Unit]
Description=Nexell update

[Service]
Type=oneshot
ExecStart=/etc/update.sh

[Install]
WantedBy=multi-user.target
