# meta-nexell-distro/recipes-qt/qt5_5.4.x/qt5-smarthome-launcher/qt5-smarthome.service
[Unit]
Description=QT5 Smarthome demo

[Service]
Type=oneshot
ExecStart=/usr/bin/qt5_smarthome_run.sh

[Install]
WantedBy=multi-user.target
