# recipes-extended/nexell-nxlog/files/32/lib/systemd/system/nxlogrotate.service
[Unit]
Description=NX log rotate daemon

[Service]
Type=idle
ExecStart=/usr/sbin/nxlogrotate

[Install]
WantedBy=multi-user.target
